# tests/fixtures/network-exposure/svc-pass.socket
[Unit]
Description=Test socket pass

[Socket]
ListenStream=127.0.0.1:5432

[Install]
WantedBy=sockets.target
